# deploy/systemd/sclib-aggregate.timer
[Unit]
Description=Fire SCLib hourly materials aggregate
Documentation=file:///opt/SCLib_JZIS/deploy/systemd/README.md

[Timer]
# Hourly at :30. Was daily originally — but the aggregator finishes
# in ~10 s on 30 k papers, and waiting up to 24 h made the materials
# count lag noticeably behind paper / chunk counts on the dashboard.
# Hourly gets the materials table within ~30 min of any new arXiv
# paper while still costing essentially nothing.
#
# :30 is chosen so we land mid-hour, between the 瓦力 ingest tick
# (typically near :00 / :05) and the API's hourly stats-refresh
# (every 3600 s from container start, drifts by minutes). No
# contention even if cadences shift.
OnCalendar=*:30:00
RandomizedDelaySec=60s
Persistent=true
Unit=sclib-aggregate.service
AccuracySec=30s

[Install]
WantedBy=timers.target
