# deploy/systemd/sclib-ingest.timer
[Unit]
Description=Fire SCLib hourly incremental ingest
# If the VPS was down across one or more scheduled firings, make up
# for exactly one missed run on next boot. More than one would
# just hammer arXiv pointlessly — the harvest cursor already tracks
# where we left off.
Documentation=file:///opt/SCLib_JZIS/deploy/systemd/README.md

[Timer]
# Every hour at :10. Off the top of the hour to avoid the default cron
# stampede (most things run at :00) and to give Postgres a quiet
# window in case daily cron is still wrapping up.
OnCalendar=*-*-* *:10:00
# Randomize by up to 60s so multiple nodes (if we ever add one) don't
# hit arXiv's OAI-PMH endpoint simultaneously.
RandomizedDelaySec=60s
# Catch up on one missed firing if the VPS was off at :10.
Persistent=true
# Unit to activate when the timer fires. Matches the .service file name.
Unit=sclib-ingest.service
# Accuracy below 1s is overkill here; default 1min is fine.
AccuracySec=30s

[Install]
WantedBy=timers.target
