# deploy/systemd/sclib-ingest.service
[Unit]
Description=SCLib hourly incremental ingest
# Don't kick off until Docker is ready — the script spins up a
# one-shot `ingestion` container via `docker compose run`.
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
# Scheduled work runs as root for docker CLI access. The underlying
# pipeline code runs as UID 1001 (sclib) inside the container, so the
# host-side root privilege is narrow.
User=root
Group=root
# Repo path on VPS2. Override via /etc/systemd/system/sclib-ingest.service.d/override.conf
# if the repo ever moves.
Environment=SCLIB_ROOT=/opt/SCLib_JZIS
Environment=SCLIB_LOG_DIR=/var/log/sclib
Environment=SCLIB_STATE_DIR=/var/lib/sclib

# The orchestrator uses flock internally and is safe to exec directly.
ExecStart=/opt/SCLib_JZIS/scripts/sclib-hourly-ingest.sh

# Resource guardrails. An hourly run shouldn't eat more than ~1 GB of
# RAM on the host side (the ingestion container has its own limits).
# Kill and log if we blow past.
MemoryMax=1G
TimeoutStartSec=30min

# If the orchestrator exits non-zero, systemd marks the unit failed.
# SuccessExitStatus=99 maps "skipped due to flock" to a clean status
# so `systemctl status` doesn't turn red when a previous run overruns.
SuccessExitStatus=99

# systemd journal captures stdout/stderr; the orchestrator also writes
# to /var/log/sclib/hourly.log for a human-friendly tail.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
